# github-trending/com.lihongwei.sync-stats.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lihongwei.sync-stats</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/huangpeng/Desktop/lihongwei-cn/github-trending/sync_stats.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/huangpeng/Desktop/lihongwei-cn/github-trending/sync_stats.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/huangpeng/Desktop/lihongwei-cn/github-trending/sync_stats_error.log</string>
</dict>
</plist>
